# problem03/infra/roles/app-install/templates/plop.service.j2
[Unit]
Description=Genius app

[Service]
ExecStart=/root/plop.py
Restart=always

[Install]
WantedBy=multi-user.target
